# k8s/ribc-scheduler/busscheduler.timer
[Unit]
Description=Run RIBC Bus and Driver Scheduler weekly

[Timer]
OnCalendar=Sun 01:00
Persistent=true

[Install]
WantedBy=timers.target